# root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/usr/bin/with-contenv bash
exec &>> >(tee -a /var/log/svc-cron)

if [ -z ${NOCRON+x} ]; then
  echo "[Cron] Starting crond"  
  exec /usr/sbin/crond -f -S -l 5 -c /etc/crontabs
fi
