# root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/usr/bin/with-contenv bash

if [ -z ${NOCRON+x} ]; then
  echo "[Cron] Starting crond" | tee /var/log/svc-cron-init
  exec /usr/sbin/crond -f -S -l 5 -c /etc/crontabs
fi